# config/cron.d/example
* * * * * root echo testing > /cron.log
